# cw-control/src/screenconnect-onlogin.plist
<?xml version="1.0" encoding="utf-8"?>
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>screenconnect-<THUMB>-onlogin</string>
    <key>ProgramArguments</key>
    <array>
      <string>/opt/screenconnect-<THUMB>.app/Contents/MacOS/ConnectWiseControl Client</string>
      <string>MAIN_ACTION=start_client</string>
      <string>MAIN_ACTION_ARGUMENT=has_service</string>
    </array>
    <key>LimitLoadToSessionType</key>
    <string>Aqua</string>
    <key>KeepAlive</key>
    <dict>
      <key>SuccessfulExit</key>
      <false />
    </dict>
    <key>EnvironmentVariables</key>
    <dict>
      <key>LC_CTYPE</key>
      <string>en_US.UTF-8</string>
      <key>LANG</key>
      <string>C</string>
      <key>SCREENCONNECT_IS_PRELOGIN</key>
      <string>False</string>
    </dict>
  </dict>
</plist>